# config/gracc-gold.cron
MAILTO=rynge@isi.edu

15 * * * * gold /usr/bin/gracc-gold -s 120
